# n1/files/etc/init.d/bypass-ipv6-fix
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /bin/sh -c '
        LOGFILE=/tmp/uci-defaults-log.txt
        LAN_DEV=br-lan
        for round in 1 2 3 4 5; do
            if ip -6 route show default | grep -q "^default "; then
                echo "[bypass-ipv6-fix] default route already exists" >> "$LOGFILE"
                exit 0
            fi
            sleep 15
            IPV6_GW=$(ip -6 neigh show dev "$LAN_DEV" | awk "/router/ && /fe80::/ {print \$1; exit}")
            if [ -n "$IPV6_GW" ]; then
                if ip -6 route add default via "$IPV6_GW" dev "$LAN_DEV" 2>/dev/null; then
                    echo "[bypass-ipv6-fix] added default route via $IPV6_GW on $LAN_DEV (round $round)" >> "$LOGFILE"
                    exit 0
                fi
            fi
        done
        echo "[bypass-ipv6-fix] failed to recover IPv6 default route after retries" >> "$LOGFILE"
        exit 0
    '
    procd_set_param respawn 0 0 0
    procd_close_instance
}
