# files/etc/uci-defaults/99-custom-docker
#!/bin/sh

# Docker 配置目录
mkdir -p /etc/docker
mkdir -p /opt/docker-apps

# Docker 默认配置：保留 IPv6 能力，但不主动启用 Docker 专属 IPv6
cat > /etc/docker/daemon.json <<'EOC'
{
  "iptables": true,
  "data-root": "/opt/docker-apps",
  "log-driver": "json-file",
  "log-opts": {
    "max-size": "10m",
    "max-file": "3"
  }
}
EOC

# 配置 Docker 防火墙区域与转发，确保 LAN 可访问容器映射端口
uci -q delete firewall.docker
for idx in $(uci show firewall | grep "=forwarding" | cut -d[ -f2 | cut -d] -f1 | sort -rn); do
    src=$(uci get firewall.@forwarding[$idx].src 2>/dev/null)
    dest=$(uci get firewall.@forwarding[$idx].dest 2>/dev/null)
    if [ "$src" = "docker" ] || [ "$dest" = "docker" ]; then
        uci -q delete firewall.@forwarding[$idx]
    fi
done

cat <<'EOF' >> /etc/config/firewall

config zone 'docker'
  option name 'docker'
  option input 'ACCEPT'
  option output 'ACCEPT'
  option forward 'ACCEPT'
  list subnet '172.16.0.0/12'

config forwarding
  option src 'lan'
  option dest 'docker'

config forwarding
  option src 'docker'
  option dest 'lan'

config forwarding
  option src 'docker'
  option dest 'wan'
EOF

uci commit firewall

# 启用 Docker 服务（如果已编进固件）
if [ -x /etc/init.d/dockerd ]; then
    /etc/init.d/dockerd enable
fi

# 重载防火墙
/etc/init.d/firewall restart >/dev/null 2>&1 || true

exit 0
